# runtime/launchagent/com.companion.ingest-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.companion.ingest-watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd ~/Documents/Companion && npm run ingest:watch</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>~/Documents/Companion/90_logs/ingest-watch.log</string>

  <key>StandardErrorPath</key>
  <string>~/Documents/Companion/90_logs/ingest-watch.log</string>
</dict>
</plist>
