# runtime/launchagent/com.companion.ingestwatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.companion.ingestwatch</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd "__COMPANION_ROOT__" &amp;&amp; npm run ingest:watch</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__COMPANION_ROOT__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__COMPANION_ROOT__/90_logs/ingest-watch.log</string>
    <key>StandardErrorPath</key>
    <string>__COMPANION_ROOT__/90_logs/ingest-watch.log</string>
  </dict>
</plist>

